# etc/sv/opendkim/run
#!/bin/sh
exec 2>&1

DAEMON=/usr/sbin/opendkim
[ -x "$DAEMON" ] || exit 0

[ -r /etc/opendkim.conf ] && CONF="-x /etc/opendkim.conf"

exec $DAEMON -f $CONF
