# etc/sv/smartd/run
#!/bin/sh -eu
exec /usr/sbin/smartd --no-fork
