# etc/sv/smartd/run
#!/bin/sh
exec /usr/sbin/smartd --no-fork
