# etc/sv/bootlogd/run
#!/bin/sh
exec 2>&1

[ -x /sbin/bootlogd ] || exit 0

# stop any bootlog daemons running outside of runsvdir
killall --quiet --exact bootlogd

exec /sbin/bootlogd -d -c -l /var/log/boot.log
