# etc/sv/cron/run
#!/bin/sh
exec 2>&1

DAEMON=/usr/sbin/cron
[ -x "$DAEMON" ] || exit

sv start socklog >/dev/null

exec "$DAEMON" -f -l -L 1
