# etc/sv/collectd/run
#!/bin/sh
exec 2>&1

NAME=collectd
DAEMON=$(command -v "$NAME")
[ -x "$DAEMON" ] || exit 0

exec "$DAEMON" -f
