# etc/sv/amazon-ssm-agent/run
#!/bin/sh

DAEMON=$(command -v amazon-ssm-agent)
[ -x "$DAEMON" ] || exit

exec $DAEMON >/dev/null 2>/dev/null
