# etc/sv/amazon-ssm-agent/run
#!/bin/sh
exec 2>&1

DAEMON=$(command -v amazon-ssm-agent)
[ -x "$DAEMON" ] || exit

exec $DAEMON
