# etc/sv/autofs/run
#!/bin/sh
exec 2>&1
[ -x /usr/sbin/automount ] || exit 0
exec /usr/sbin/automount --foreground --verbose 2>&1
