# etc/sv/autofs/run
#!/bin/sh -eu
exec /usr/sbin/automount --foreground
