# etc/sv/gpg-agent/run
#!/bin/sh
exec 2>&1
DAEMON=$(which gpg-agent)
[ -x "$DAEMON" ] || exit 0
gpgconf --kill gpg-agent
exec $DAEMON --supervised
