# etc/sv/dhclient/run
#!/bin/sh
exec 2>&1

[ -x /sbin/dhclient ] || exit 0

[ -r ./conf ] && . ./conf

exec /sbin/dhclient ${OPTS:=-d -4 -w -I}
