# etc/sv/networking/finish
#!/bin/sh
service=${PWD##*/}

# run this only once
sv stop $service
