# etc/sv/incron/run
#!/bin/sh
# The inotify cron daemon (incrond) monitors filesystem events
# and executes commands defined in system and user tables.
# (cloux@rote.ch)
exec 2>&1

NAME=incrond
DAEMON=$(command -v $NAME)
[ -x "$DAEMON" ] || exit

# WARNING:
# If incrond can't connect to syslog service by openlog() call
# it will keep polluting /dev/console.
sv start socklog >/dev/null || sleep 1
socklog-check 2>/dev/null || sleep 1

exec "$DAEMON" --foreground -f /etc/incron.conf >/dev/null 2>/dev/null
