# etc/sv/incron/run
#!/bin/sh
exec 2>&1

INCROND=/usr/sbin/incrond
INCROND_CONF=/etc/incron.conf

[ -x "$INCROND" ] || exit 0

# incron should run even without logging, so don't exit here if socklog fails:
sv start socklog-unix >/dev/null
# incron jobs keep polluting /dev/console, do not log that
sv stop bootlogd >/dev/null
killall --quiet --exact bootlogd

exec "$INCROND" --foreground -f "$INCROND_CONF"
