# etc/sv/rpc.mountd/run
#!/bin/sh
exec 2>&1
[ -x /usr/sbin/rpc.mountd ] || exit
sv start rpc.nfsd || exit
exec /usr/sbin/rpc.mountd --foreground --no-nfs-version 2 --no-nfs-version 3
