# etc/sv/udevd/run
#!/bin/sh
exec 2>&1

[ -x /sbin/udevd ] || exit

# udevd could be started unsupervised in the early init.
# kill the old daemon first
/sbin/udevadm control --exit 2>/dev/null

exec /sbin/udevd
