# etc/sv/pure-ftpd/run
#!/bin/sh
exec 2>&1
exec pure-ftpd -a 50 -E -l pam -u 100
