# etc/sv/postgresql/run
#!/bin/sh
# PostgreSQL Cluster
# see 'man pg_ctlcluster' and 'man pg_ctl' for details
#
# (cloux@rote.ch)
exec 2>&1

DAEMON="$(command -v pg_ctlcluster)"
[ -x "$DAEMON" ] || exit

# ensure sysvinit compatibility - replace init.d script with symlink to sv
NAME=${PWD##*/}
if [ -f /etc/init.d/$NAME ] && [ ! -h /etc/init.d/$NAME ]; then
	mv -f /etc/init.d/$NAME /etc/init.d/${NAME}.sysv
	ln -sf "$(command -v sv)" /etc/init.d/$NAME
fi

NAME=
OPTIONS=
VERSION=
PGBINROOT="/usr/lib/postgresql/"
[ -r ./conf ] && . ./conf
[ -d "$PGBINROOT" ] || exit 1
[ -z "$VERSION" ] && VERSION=$(ls -1v $PGBINROOT | tail -n 1)
[ -z "$NAME" ] && NAME=$(ls -1tr /etc/postgresql/$VERSION | tail -n 1)

install -d -m 2775 -o postgres -g postgres -d /run/postgresql
exec chpst -u postgres $DAEMON --foreground --skip-systemctl-redirect $VERSION $NAME start -- $OPTIONS
