# etc/sv/nginx/run
#!/bin/sh
exec 2>&1
# nginx only provides a config file option to run in the foreground:
exec /usr/sbin/nginx -c /etc/nginx/nginx.conf -g 'daemon off;'
