# etc/sv/hiawatha/run
#!/bin/sh -e
#
# Hiawatha Webserver (https://www.hiawatha-webserver.org)

exec 2>&1
[ -x "$(which hiawatha 2>/dev/null)" ] || exit 1
exec hiawatha -d -c /etc/hiawatha >/dev/null
